# notebooks/launchd/com.kelvinn.healthdata.turso-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kelvinn.healthdata.turso-sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>set -a; [ -f "{{ENV_FILE}}" ] &amp;&amp; source "{{ENV_FILE}}"; set +a; "{{PYTHON_PATH}}" "{{SCRIPT_PATH}}" sync</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>StandardOutPath</key>
  <string>{{STDOUT_LOG}}</string>

  <key>StandardErrorPath</key>
  <string>{{STDERR_LOG}}</string>
</dict>
</plist>
